# os-config/rootfs_overlay/etc/inittab
::sysinit:/etc/init.d/rcS
tty1::respawn:env PATH=/usr/bin:/usr/sbin:/bin:/sbin TERM=linux /usr/bin/fix-automation
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/init.d/rcK
::restart:/sbin/init