# os-config/rootfs_overlay/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:env TERM=linux /usr/bin/fix-automation
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/init.d/rcK
::restart:/sbin/init