# os-config/rootfs_overlay/etc/inittab
::sysinit:/etc/init.d/rcS
tty1::respawn:env TERM=linux setsid /usr/bin/fix-automation
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/init.d/rcK
::restart:/sbin/init